# templates/profile.d/01-createtmp.sh
#! /usr/bin/env bash

TMPDIR="${HOME}/.local/tmp"
TEMP="${TMPDIR}"
TMP="${TMPDIR}"
mktemp