# templates/systemd/goodnight.timer.j2
[Unit]
Description=Schedule goodnight execution dail at {{ goodnight_time }}

[Timer]
OnCalendar=*-*-* {{ godnight_time }}
Unit=goodnight.service

[Install]
WantedBy=default.target